# supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /thingin_Adapter/thingin_recommender_uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx